# crontab
# Pipeline ETL - Receita Federal
# Executa todos os dias as 06:00 AM
# Verifica o ambiente a cada 6 horas
# Remove logs com mais de 7 dias
@"
0 6 * * * root cd /app && /usr/local/bin/python etl/main.py >> logs/cron.log 2>&1
0 */6 * * * root cd /app && /usr/local/bin/python 0_check_ambiente.py >> logs/health_check.log 2>&1
0 0 * * 0 root find /app/logs -name "*.log" -mtime +7 -delete
"@ | Set-Content -Path "crontab" -Encoding UTF8

