# crontab
# Pipeline ETL - Receita Federal# Executa todos os dias às 06:00 AM
0 6 * * * cd /app && python etl/main.py >> logs/cron.log 2>&1

# Verificação de saúde do sistema - a cada hora
0 * * * * cd /app && python 0_check_ambiente.py >> logs/health_check.log 2>&1